# crontab
0 2 * * * /home/ec2-user/mediawiki-docker/backup.sh